# build/never-todo-backend.service
[Unit]
Description=mpgo

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/never-todo/main start -c /etc/never-todo.ini

[Install]
WantedBy=multi-user.target